# reenable-gha.timer
[Unit]
Description=Run reenable-gha monthly

[Timer]
OnCalendar=monthly
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
